# systemd/lanmon-agent.service
[Unit]
Description=LAN Discovery Agent (LANBeacon)
Documentation=https://github.com/lanmon/lanmon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lanmon
Group=lanmon
ExecStart=/usr/local/bin/lanmon agent
Restart=on-failure
RestartSec=10s

# Security hardening
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/lanmon /run/lanmon
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
